# hot-corners.service
[Unit]
Description=hot-corners

[Service]
ExecStart=
Environment="DISPLAY="
Environment="XAUTHORITY="

[Install]
WantedBy=multi-user.target
